# deploy/systemd/sovereign-node-disk-check.timer
[Unit]
Description=Daily disk space check for Sovereign Node

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
