# deploy/systemd/sovereign-node-docker-prune.service
[Unit]
Description=Sovereign Node Docker image prune

[Service]
Type=oneshot
ExecStart=/usr/bin/docker image prune -f
